# fledger-experiments/playbook/flsignal.service
[Unit]
Description=Fledger signaling server
After=network.target

[Service]
Type=simple
ExecStart=/root/flsignal -v --max-list-len 25
User=root
Restart=no
StandardOutput=append:/var/log/flsignal

[Install]
WantedBy=multi-user.target
